# devices/linksys_mx4300/files/etc/uci-defaults/999-Local_config
# 999-Local_config - LazyWoodpecker local configuration

uci -q batch << EOF

# Timezone
set system.@system[0].timezone='EST5EDT,M3.2.0,M11.1.0'
set system.@system[0].zonename='America/New York'

EOF

exit 0
